# src/elasticsearch/etc/supervisor.d/elasticsearch.conf
[program:elasticsearch]
command=/elasticsearch/bin/elasticsearch -Des.default.config=/elasticsearch/config/elasticsearch.yml
priority=5
