# src/kibana/etc/supervisor.d/kibana.conf
[program:kibana]
command=hypnotoad -f /opt/kibana-authorization/kbnauth/script/kbnauth
priority=5
